# src/installer/dev/config/hooks/normal/0200-autologin.hook.chroot
#!/bin/bash
apt-get remove --purge -y sddm 2>/dev/null || true

mkdir -p /etc/systemd/system/getty@tty1.service.d
cat > /etc/systemd/system/getty@tty1.service.d/autologin.conf << CONF
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin user --noclear %I $TERM
CONF

mkdir -p /etc/init.d
cat > /etc/init.d/tty1-autologin << CONF
#!/sbin/openrc-run
command="/sbin/agetty"
command_args="--autologin user --noclear tty1 linux"
command_background=true
pidfile="/run/tty1-autologin.pid"
depend() {
    need localmount
}
CONF
chmod +x /etc/init.d/tty1-autologin
rc-update add tty1-autologin default 2>/dev/null || true

cat > /etc/skel/.bash_profile << CONF
if [ -z "\$DISPLAY" ] && [ "\$(tty)" = "/dev/tty1" ]; then
    exec dbus-run-session startplasma-wayland
fi
CONF

cat > /root/.bash_profile << CONF
if [ -z "\$DISPLAY" ] && [ "\$(tty)" = "/dev/tty1" ]; then
    exec dbus-run-session startplasma-wayland
fi
CONF
